# launcher/factorio@.service
[Unit]
Description=Factorio Server %i

[Service]
User=factorio
Group=factorio
WorkingDirectory=/opt/factorio
Type=simple
SyslogIdentifier=factorio-%i
KillSignal=SIGINT
ExecStart=/opt/factorio/launcher/factorio.sh env-%i

RestartSec=10
Restart=on-failure

[Install]
WantedBy=multi-user.target
